# wireguard-api/root/etc/s6-overlay/s6-rc.d/wg-api/run
#!/command/with-contenv bash
cd /app
exec node server.js
